# systemd/system/mdmonitor.service.d/sandbox.conf
[Service]
MemoryDenyWriteExecute=yes
IPAddressDeny=any
PrivateNetwork=yes
PrivateUsers=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
