# systemd/system/dehydrated.service.d/apache2.conf
[Service]
ExecStartPost=!/bin/systemctl reload apache2
