# systemd/system/dehydrated.service.d/dovecot.conf
[Service]
ExecStartPost=!/bin/systemctl reload dovecot
