# systemd/system/dehydrated.service.d/postfix.conf
[Service]
ExecStartPost=!/bin/systemctl reload postfix@-
