# systemd/system/mail-backup@.timer
[Unit]
Description=Back up /var/mail/vmail/lucaswerkmeister.de %i

[Timer]
OnCalendar=%i
Persistent=true

[Install]
WantedBy=timers.target
