# systemd/system/mariadb.service.d/localhost.conf
[Service]
IPAddressAllow=localhost
IPAddressDeny=any
RestrictNetworkInterfaces=lo
