# systemd/system/dovecot.service.d/ProtectSystem.conf
[Service]
# vendor unit file instructs us to “enable this if your systemd is new enough to support it”; it is:
ProtectSystem=full
