# systemd/system/phpsessionclean.timer.d/override.conf
[Unit]
Description=Clean PHP session files every week

[Timer]
OnCalendar=
OnCalendar=weekly
