# systemd/system/tor@default.service.d/RuntimeMaxSec.conf
# restart Tor every 12 hours to work around memory leak :/
[Service]
RuntimeMaxSec=12h
Restart=on-failure
